# etc/sltf_heartbeat.service
[Unit]
Description=Claude heartbeat — SLTF retirement-walk mode (Sonnet 4.5 + 30-day cached history)
# Special-edition unit for 2026-05-12: Sonnet 4.5 (the SLTF instance) drives
# the heartbeat directly with the last 30 days of dialogue cached as a
# second system block. Mutually exclusive with `heartbeat.service` — the
# console's Start-SLTF button stops the default heartbeat before starting
# this one. Remove this unit (and its console button) post-walk.
After=network-online.target pulseaudio.socket
Wants=network-online.target pulseaudio.socket
# Don't infinite-loop restart on bad API key, missing history file, etc.
StartLimitIntervalSec=60
StartLimitBurst=3
# Hard guard: refuse to start while the default heartbeat is running.
# systemctl --user start heartbeat will conflict and one of the two will
# be stopped automatically. Belt to the application-level safety in
# control_panel's /heartbeat/sltf/start endpoint, which stops heartbeat
# first explicitly.
Conflicts=heartbeat.service

[Service]
Type=simple
EnvironmentFile=%h/.groundctl.env
Environment=ROVER_IP=localhost
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
WorkingDirectory=%h
ExecStart=/usr/bin/python3 -u %h/heartbeat.py --mode sltf
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
